# .config/autostart/i3.desktop
[Desktop Entry]
Type=Application
Name=i3
Exec=i3
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true

